# examples/packer/al2-ebs-mockend-systemd/payload/mockend.service
[Unit]
Description=mockend
[Service]
ExecStart=/usr/bin/java -jar /tmp/mockend-runner.jar